# conf/ovirt-server.crontab
MAILTO=root
*/5 * * * * root /usr/sbin/ovirt-update-search
